# services/pedalboard-bridge.service
[Unit]
Description=Pedalboard Bridge (JACK MIDI + Audio Switching)
After=pedalboard-modhost.service
Wants=pedalboard-jack.service pedalboard-modhost.service

[Service]
LimitRTPRIO=95
LimitMEMLOCK=infinity
Type=simple
User=laenzi
EnvironmentFile=/etc/pedalboard/env
ExecStart=/udata/pedalboard-bridge/pedalboard-bridge-rust --midi ${MIDI_MATCH} --addr 0.0.0.0:${BRIDGE_PORT} --audio /etc/pedalboard/audio-rig.yaml --modhost localhost:5555
Restart=always
RestartSec=2
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
